# script/docker/supervisord.conf.example
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:php7-fpm]
command=/usr/sbin/php-fpm7.1 -F